# distro/rootfs/overlay/etc/systemd/system/munin-sts.service
[Unit]
Description=MuninOS Speech-to-Speech Runtime
After=sound.target network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/munin-sts
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
